# supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
# Environment variables are inherited from Docker container
# Use %(ENV_VAR_NAME)s in program sections to access them
# VastAI port variables (VAST_TCP_PORT_70000, etc.) are set by VastAI and must be
# explicitly exported in the onstart script before supervisord starts

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:watcher
# Note: chown only works when socket is first created
# If socket exists from previous run, it may need manual permission fix
# The monitor service will attempt to fix permissions but requires root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:jupyter]
command=/usr/local/bin/jupyter server --port=8888 --IdentityProvider.token=%(ENV_JUPYTER_TOKEN)s --ServerApp.disable_check_xsrf=true
directory=/app
user=appuser
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/jupyter.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
# Environment variables are inherited from supervisord (which inherits from Docker)
# Only set static variables here; dynamic ones come from Docker -e flags
# Note: CSRF checks are disabled because we use token authentication which provides security
environment=USER="appuser",HOME="/home/appuser",XDG_DATA_HOME="/tmp/.local/share",XDG_CONFIG_HOME="/tmp/.config",XDG_CACHE_HOME="/tmp/.cache",JUPYTER_RUNTIME_DIR="/tmp/jupyter-runtime",JUPYTER_DATA_DIR="/tmp/jupyter-data",JUPYTER_CONFIG_DIR="/tmp/jupyter-config",TD_LOG_DIR="/tmp/tensordock-logs"

[program:python_server]
command=/usr/local/bin/python -u /app/run_modular.py
directory=/app
user=appuser
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/python_server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=USER="appuser",HOME="/home/appuser",XDG_DATA_HOME="/tmp/.local/share",XDG_CONFIG_HOME="/tmp/.config",XDG_CACHE_HOME="/tmp/.cache",TD_LOG_DIR="/tmp/tensordock-logs"

[program:turn_server]
command=/bin/bash -c "if [ \"%(ENV_START_TURN)s\" = \"true\" ]; then TURN_PORT=\"%(ENV_VAST_UDP_PORT_70001)s\"; if [ -z \"$TURN_PORT\" ]; then TURN_PORT=\"50000\"; fi; PUBLIC_IP=\"%(ENV_PUBLIC_IPADDR)s\"; if [ -z \"$PUBLIC_IP\" ] || [ \"$PUBLIC_IP\" = \"auto\" ]; then META_IP=$(curl -fsS --max-time 2 http://169.254.169.254/metadata/v1/interfaces/public/0/ipv4/address 2>/dev/null || true); if [ -n \"$META_IP\" ]; then PUBLIC_IP=\"$META_IP\"; else PUBLIC_IP=$(curl -fsS --max-time 3 https://ifconfig.co 2>/dev/null || true); fi; fi; if [ -n \"$PUBLIC_IP\" ]; then /usr/bin/turnserver -n -a --log-file=stdout --lt-cred-mech --fingerprint --no-stun --no-multicast-peers --no-cli --no-tcp --no-tls --realm=\"example.org\" --user=\"%(ENV_TURN_USERNAME)s:%(ENV_TURN_PASSWORD)s\" -p \"$TURN_PORT\" -X \"$PUBLIC_IP\"; else /usr/bin/turnserver -n -a --log-file=stdout --lt-cred-mech --fingerprint --no-stun --no-multicast-peers --no-cli --no-tcp --no-tls --realm=\"example.org\" --user=\"%(ENV_TURN_USERNAME)s:%(ENV_TURN_PASSWORD)s\" -p \"$TURN_PORT\"; fi; else exit 0; fi"
directory=/app
user=appuser
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/turn_server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
# Environment variables are inherited from supervisord (which inherits from Docker)
# VAST_UDP_PORT_70001 is the external port (identity mapping - same as internal)
# Only set static variables here; dynamic ones come from Docker -e flags
environment=USER="appuser",HOME="/home/appuser",VAST_UDP_PORT_70001="%(ENV_VAST_UDP_PORT_70001)s",PUBLIC_IPADDR="%(ENV_PUBLIC_IPADDR)s",TURN_USERNAME="%(ENV_TURN_USERNAME)s",TURN_PASSWORD="%(ENV_TURN_PASSWORD)s"

[program:monitor]
command=/usr/local/bin/python -u /app/monitor_service.py
directory=/app
user=watcher
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/monitor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
# Environment variables are inherited from supervisord (which inherits from Docker)
# VastAI port variables (VAST_TCP_PORT_70000, etc.) are passed via %(ENV_VAR_NAME)s
# These are set by VastAI and exported in the onstart script before supervisord starts
environment=USER="watcher",HOME="/home/watcher",VAST_TCP_PORT_70000="%(ENV_VAST_TCP_PORT_70000)s",VAST_UDP_PORT_70001="%(ENV_VAST_UDP_PORT_70001)s",VAST_TCP_PORT_70002="%(ENV_VAST_TCP_PORT_70002)s"

